# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting nginx..."
exec nginx
