# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Run nginx
bashio::log.info "Starting nginx..."
exec nginx