# rootfs/etc/s6-overlay/s6-rc.d/backend/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting dnsmasq-DHCP backend..."
exec /opt/bin/backend