# rootfs/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bashio
CONFIG="/etc/dnsmasq.conf"

# Run nginx
bashio::log.info "Starting nginx..."
nginx -g "daemon off;error_log /dev/stdout debug;" &

# Run dnsmasq
bashio::log.info "Starting dnsmasq..."

# Set max open file limit to speed up startup
ulimit -n 1024

exec dnsmasq -C "${CONFIG}" -z < /dev/null

# useful for debug, to launch a container doing nothing to inspect templating results:
#exec sleep 10000

