# rootfs/etc/services.d/web-backend/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting dhcp-clients-webapp-backend..."
exec /opt/bin/dhcp-clients-webapp-backend

